# dotfiles/systemd/firefox.service
# {{ header() }}

[Unit]
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/firefox

[Install]
WantedBy=sway-session.target
