# dotfiles/systemd/sway-idle.service
# {{ header() }}

[Unit]
Description= Send notification when maildir gets modified
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/bin/env _sway_idle

[Install]
WantedBy=sway-session.target
