# dotfiles/systemd/telegram-desktop.service
# {{ header() }}

[Unit]
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/telegram-desktop
Restart=always
StartLimitBurst=2
StartLimitInterval=1

[Install]
WantedBy=sway-session.target
